# camera-job.conf
description "run the camera code at startup pls"
author "Sean"

start on runlevel [2345]
exec /home/ubuntu/Downloads/Camera-Vision-2016/Camera-Vision-2016 quiet > /var/log/camera-job.log
